# talking.service
[Unit]
Description=talking chat
After=network.target

[Service]
WorkingDirectory=/opt/talking
ExecStart=/usr/bin/node server.js
EnvironmentFile=/opt/talking/.env
Restart=always
RestartSec=3
User=talking
Group=talking
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/talking/data
ProtectHome=true

[Install]
WantedBy=multi-user.target
